# src/maddy-reload.path
[Path]
PathChanged=/etc/maddy-config-grug
PathChanged=/etc/server.env

[Install]
WantedBy=multi-user.target
